# provisioning/crontab
*/3 * * * * /bin/bash -c /app/run-sync.sh
